# lib/systemd/system/updsched.timer
[Unit]
Description=Schedule upgrades

[Timer]
OnUnitActiveSec=1d
OnBootSec=15m

[Install]
WantedBy=timers.target
