# usr/lib/systemd/system/updsched.service
[Unit]
Description=Upgrade system

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/bin/updsched.sh
